# sender.timer
[Unit]
Description=Sender Timer Delay

[Timer]
OnBootSec=1min
Unit=sender.service

[Install]
WantedBy=multi-user.target
